# docker/development/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:redis]
command=redis-server --appendonly yes --dir /var/lib/redis
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis.log
autorestart=true
user=root
priority=100

[program:qdrant]
command=qdrant --config-path /app/qdrant-config.yaml
stdout_logfile=/var/log/supervisor/qdrant.log
stderr_logfile=/var/log/supervisor/qdrant.log
autorestart=true
user=root
environment=QDRANT__LOG_LEVEL=INFO,QDRANT__SERVICE__HTTP_PORT=6333
priority=200

[program:api]
command=sh -c "uvicorn api.main\:app \$UVICORN\_FLAGS"
directory=/app
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=2
stderr_logfile_backups=2
autorestart=true
user=root
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1",CELERY_BROKER_URL="redis://localhost:6379/0",CELERY_RESULT_BACKEND="redis://localhost:6379/0",QDRANT_HOST="${QDRANT_HOST}",QDRANT_PORT="${QDRANT_PORT}",UVICORN_FLAGS="${UVICORN_FLAGS}"
priority=300
startretries=5
startsecs=10

[program:celery]
command=celery -A api.main.celery_app worker -l info
directory=/app
stdout_logfile=/var/log/supervisor/celery.log
stderr_logfile=/var/log/supervisor/celery.log
autorestart=true
user=root
environment=PYTHONPATH=/app,CELERY_BROKER_URL=redis://localhost:6379/0,CELERY_RESULT_BACKEND=redis://localhost:6379/0,QDRANT_HOST="${QDRANT_HOST}",QDRANT_PORT="${QDRANT_PORT}"
priority=400

[program:streamlit]
command=streamlit run ui/Home.py --server.address 0.0.0.0 --server.port 8501 --server.enableCORS false --server.enableXsrfProtection false
directory=/app
stdout_logfile=/var/log/supervisor/streamlit.log
stderr_logfile=/var/log/supervisor/streamlit.log
autorestart=true
user=root
environment=PYTHONPATH=/app,API_BASE_URL=http://localhost:18000,STREAMLIT_SERVER_ENABLE_CORS=false,STREAMLIT_SERVER_ENABLE_XSRF_PROTECTION=false
priority=500